# systemd/backup.timer

[Unit]
Description=Ejecutar Copia de Respaldo

[Timer]
OnCalendar=*-*-* 00:13:00
Persistent=true

[Install]
WantedBy=multi-user.target
